# shard-agent/claude-shard-agent.service
[Unit]
Description=Claude Shard Agent
Documentation=https://github.com/myctobot/shard-agent
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/shard-agent
# Adjust this path if node is installed elsewhere (check with: which node)
ExecStart=/root/.nvm/versions/node/v24.12.0/bin/node /opt/shard-agent/src/index.js
Restart=on-failure
RestartSec=10

# Environment
Environment=NODE_ENV=production
EnvironmentFile=/opt/shard-agent/.env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=claude-shard-agent

# Security
NoNewPrivileges=false

[Install]
WantedBy=multi-user.target
